# debian/postrm
#!/bin/sh

#DEBHELPER#

set -e

rm -f /etc/gunicorn.d/musicdb
gunicorn-debian reload musicdb

rm -f /etc/nginx/sites-enabled/musicdb
invoke-rc.d nginx reload
